# eaas-base/files/service/apache2/run
#!/bin/sh

exec /usr/sbin/apachectl -DFOREGROUND >> /var/log/apache.log 2>&1
