# eaas-base/files/service/10_mongodb/run
#!/bin/bash
su - bwfla
exec /usr/bin/mongod --dbpath /home/bwfla/server-data/db --smallfiles
