# eaas-base/files/service/10_mongodb/run
#!/bin/sh

dbpath='/eaas/server-data/db'

chown -R bwfla:bwfla "${dbpath}"
exec /sbin/setuser bwfla /usr/bin/mongod \
	--dbpath "${dbpath}" \
	--smallfiles
