# eaas-base/files/service/10_mongodb/run
#!/bin/bash
su - bwfla
mkdir -p /home/bwfla/server-data/db
exec /usr/bin/mongod --dbpath /eaas/server-data/db --smallfiles
