# eaas-base/files/service/10_mongodb/run
#!/bin/bash
exec /usr/bin/mongod --dbpath /eaas/server-data/db --smallfiles
